# loopback/dell5285-camera-loopback.service
[Unit]
Description=Dell Latitude 5285 camera loopback (libcamera -> v4l2loopback)
After=graphical-session.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/.local/bin/dell5285-camera-loopback
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
